# raspbian/usbnetwork.service
[Unit]
Description=Enable USB network interface

[Service]
Type=oneshot
ExecStart=/opt/sbin/usbnetworkstart.sh

[Install]
WantedBy=multi-user.target
